# etc/periodic/daily/810.vulture_logrotate
#!/bin/sh
#
# This file is part of VultureOS
#

/usr/local/sbin/logrotate -f /usr/local/etc/logrotate.conf
